# systemd/raspi_fanspeed.service
[Unit]
Description=adjustable fanspeed with temperature monitoring
Requires=pigpiod.service
After=network-online.target

[Service]
Type=simple
User=root
ExecStart=$RPI_FANSPEED_COMMAND_LINE

[Install]
Alias=fanspeed.service
WantedBy=multi-user.target
